# deploy/roles/common/templates/gunicorn.socket
[Unit]
Description={{app_name}} gunicorn socket

[Socket]
ListenStream=/run/{{ gunicorn_name }}.sock

[Install]
WantedBy=sockets.target
